# fitness-app-watch.service
[Unit]
Description=Fitness App Auto-Deploy Watcher
After=docker.service
Wants=docker.service
# Don't rate-limit restarts — the script handles its own backoff
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/etc/Fitness-App
ExecStart=/etc/Fitness-App/deploy.sh --watch 60
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fitness-app-watch

# Give deploys time to finish before killing
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
